# vm-systemd/cups.service
.include /lib/systemd/system/cups.service
[Unit]
ConditionPathExists=/var/run/qubes-service/cups
