# vm-systemd/ntpd.service
.include /lib/systemd/system/ntpd.service
[Unit]
ConditionPathExists=/var/run/qubes-service/ntpd
