# linux/systemd/qubes-vm@.service
[Unit]
Description=Start Qubes VM %i
After=qubes-netvm.service
Before=plymouth-quit.service

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/usr/bin/qvm-start --no-guid %i
Group=qubes
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
