# linux/systemd/qubes-reload-firewall@.timer
[Unit]
Description=Reload VM firewall to expire its temporary rules

[Timer]
OnActiveSec=1m
OnUnitActiveSec=1m
